# template/template.service
[Unit]
Description={description}

[Service]
Type=exec
Restart=on-failure
ExecStart={command}  {script}

[Install]
WantedBy=multi-user.target
